# misc/flask.ini
[program:flask]
directory=%(ENV_HOME)s/simple-ratings-server
command=%(ENV_HOME)s/simple-ratings-server/.venv/bin/uwsgi uwsgi.ini